# build/pkg/org.mochi-os.world.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.mochi-os.world</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/sbin/mochi-world</string>
		<string>-f</string>
		<string>/etc/mochi/world.conf</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>UserName</key>
	<string>_mochi</string>
	<key>GroupName</key>
	<string>_mochi</string>
	<key>StandardErrorPath</key>
	<string>/var/log/mochi-world.log</string>
</dict>
</plist>
